# service/example/finish
#!/bin/sh
exec 2>&1
svcname=$(basename $(pwd))

if [[ $1 -eq 0 || $2 -eq 15 ]]; then
	echo $svcname is stopped
else
	echo $svcname is crashed
fi
